# ansible/roles/pds/templates/pds-health-check.service.j2
[Unit]
Description=PDS Health Check — restart on failure
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pds-health-check.sh
